# openta/conf/rqworker.conf
[program:beta_rq_worker]
command=/Users/ostlund/WWW/beta/django/env/bin/python manage.py rqworker
directory=/Users/ostlund/WWW/beta/django/backend
user=ostlund
autostart=true
autorestart=true
